# deploy/omnideck-hub.service
[Unit]
Description=OmniDeck Hub
After=network.target

[Service]
Type=simple
User=%%USER%%
WorkingDirectory=%%INSTALL_DIR%%
ExecStart=/usr/bin/node %%INSTALL_DIR%%/hub/dist/index.js
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=OMNIDECK_CONFIG_DIR=%%HOME%%/.omnideck/config
Environment=OMNIDECK_SECRETS_PATH=%%HOME%%/.omnideck/secrets.yaml
Environment=OMNIDECK_DATA_DIR=%%HOME%%/.omnideck
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
